# services/webhook_server/launchagents/com.seren.cloudflared-tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.seren.cloudflared-tunnel</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/cloudflared</string>
    <string>tunnel</string>
    <string>--config</string>
    <string>/Users/brianhellemn/.cloudflared/config.yml</string>
    <string>run</string>
    <string>8f70cf6f-551f-4fb7-b1f9-0349f2f688f3</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>WorkingDirectory</key>
  <string>/Users/brianhellemn</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/brianhellemn</string>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/brianhellemn/Library/Logs/cloudflared-tunnel.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/brianhellemn/Library/Logs/cloudflared-tunnel.err.log</string>
</dict>
</plist>
