# webhook-server/launchagents/com.seren.webhook-server.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.seren.webhook-server</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/brianhellemn/Projects/github-production/webhook-server/notion_event_subscription_webhook_server_v4_enhanced.py</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>WorkingDirectory</key>
  <string>/Users/brianhellemn/Projects/github-production/webhook-server</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOST</key>
    <string>0.0.0.0</string>
    <key>PORT</key>
    <string>5001</string>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/brianhellemn/Library/Logs/webhook-server.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/brianhellemn/Library/Logs/webhook-server.err.log</string>
</dict>
</plist>
